# crontab
3,13,23,33,43,53 *      * * *   pi  python /home/pi/cwb-cache/cache_current.py
